# flag-mover.service
[Unit]
Description=To move the flag automatically
[Service]
Type=simple
ExecStart=/path/to/moverRunner.sh &
[Install]
WantedBy=multi-user.target